# sipXprovision/bin/sipxprovision.service.in
[Unit]
Description=sipxprovision
Documentation=http://sipfoundry.org
After=network.target syslog.target

[Service]
Type=oneshot
ExecStart=@SIPX_SERVICEDIR@/sipxprovision start
RemainAfterExit=true
ExecStop=@SIPX_SERVICEDIR@/sipxprovision stop
StandardOutput=journal

[Install]
WantedBy=multi-user.target