# sipXconfig/bin/mongodb-sync.service.in
[Unit]
Description=Mongodb sync
Documentation=http://sipfoundry.org
After=network.target syslog.target

[Service]
Type=oneshot
ExecStart=@SIPX_SERVICEDIR@/mongodb-sync start
RemainAfterExit=true
ExecStop=@SIPX_SERVICEDIR@/mongodb-sync stop
StandardOutput=journal

[Install]
WantedBy=multi-user.target