# sipXconfig/bin/mongo-local-arbiter.service.in
[Unit]
Description=Mongo local arbiter
Documentation=http://sipfoundry.org
After=network.target syslog.target

[Service]
Type=oneshot
ExecStart=@SIPX_SERVICEDIR@/mongo-local-arbiter start
RemainAfterExit=true
ExecStop=@SIPX_SERVICEDIR@/mongo-local-arbiter stop
StandardOutput=journal

[Install]
WantedBy=multi-user.target