# sipXconfig/bin/sipxbridge.service.in
[Unit]
Description=sipxbridge
Documentation=http://sipfoundry.org
After=network.target syslog.target

[Service]
# see man systemd.service
Type=oneshot
ExecStart=@SIPX_SERVICEDIR@/sipxbridge start
RemainAfterExit=true
ExecStop=@SIPX_SERVICEDIR@/sipxbridge stop
StandardOutput=journal

[Install]
WantedBy=multi-user.target