# sipXconfig/bin/sipxconfig.service.in
[Unit]
Description=sipxconfig
Documentation=http://sipfoundry.org
After=network.target syslog.target

[Service]
# see man systemd.service
Type=oneshot
ExecStart=@SIPX_SERVICEDIR@/sipxconfig start
RemainAfterExit=true
ExecStop=@SIPX_SERVICEDIR@/sipxconfig stop
StandardOutput=journal

[Install]
WantedBy=multi-user.target